# config/php/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command=php-fpm
numprocs=1
autostart=true
autorestart=true
stderr_logfile=/dev/stdout
stdout_logfile=/dev/stdout
priority=100
